# score-manager.service
[Unit]
Description=Score Manager Web Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/dashinaru/Git/score-manager
ExecStart=/home/dashinaru/Git/score-manager/start_server.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=score-manager

[Install]
WantedBy=default.target
